# Scripts/crontab
SHELL=/bin/bash
# MAILTO=hjstein@gmail.com
MAILTO=hjstein

*/5 * * * * $HOME/bin/run-update >>$HOME/Updatelog/coronavirus-data-`/bin/date +\%Y\%m\%d-\%H:\%M:\%S`.log 2>&1
